# zway/supervisor/balena_entry.conf
[program:balena_entry]
priority=1
directory=/usr/bin
command=/usr/bin/entry.sh
user=root
process_name=balena_entry
startsecs=0
stopsignal=INT
autorestart=false
autostart=true