# zway/supervisor/user_start.conf
[program:user_start]
priority=1
directory=/opt
command=/opt/start.sh
user=root
process_name=user_start
startsecs=0
stopsignal=INT
autorestart=false
autostart=true